# ansible/roles/systemd_services/templates/magni_navigation.service.j2
[Unit]
Description=Magni SLAM and path planning navigation
After=magni_hardware.service
Wants=magni_hardware.service

[Service]
Type=simple
User={{ target_user }}
Group={{ target_user }}
WorkingDirectory={{ ros2_ws }}
Environment=HOME={{ target_user_home }}
Environment=RMW_IMPLEMENTATION={{ rmw_implementation }}
ExecStart=/bin/bash -c "source /opt/ros/{{ ros_distro }}/setup.bash && source {{ ros2_ws }}/install/setup.bash && ros2 launch robot_bringup navigation.launch.py"
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
