# ansible/roles/systemd_services/templates/navigation.service.j2
[Unit]
Description=Magni Navigation (SLAM, Move Smooth, Routing)
After=magni-hardware.service
Requires=magni-hardware.service

[Service]
Type=simple
User={{ target_user }}
Group={{ target_user }}
WorkingDirectory={{ ros2_ws }}
Environment=HOME={{ target_user_home }}
Environment=RMW_IMPLEMENTATION={{ rmw_implementation }}
ExecStart=/bin/bash -c "source /opt/ros/{{ ros_distro }}/setup.bash && source {{ ros2_ws }}/install/setup.bash && ros2 launch robot_bringup navigation.launch.py"
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
