# ansible/roles/systemd_services/templates/teleop.service.j2
[Unit]
Description=Magni Teleop and Control Multiplexer
After=magni-hardware.service
BindsTo=magni-hardware.service

[Service]
Type=simple
User={{ target_user }}
Group={{ target_user }}
WorkingDirectory={{ ros2_ws }}
Environment=HOME={{ target_user_home }}
ExecStart=/usr/local/bin/magni-ros2-run.sh robot_bringup teleop.launch.py
Restart=on-failure
RestartSec=2
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
